# gentoo/update_gentoo_every_night/update.timer
[Unit]
Description=Update gentoo every day

[Timer]
OnCalendar=*-*-* 02:00:00


[Install]
WantedBy=timers.target
